# ansible/roles/arch/systemd/files/user/battery-notifications.service
[Unit]
Description=Battery Notifications

[Service]
Type=simple
ExecStart=/home/markl/.scripts/battery-notifier.sh
