# ansible/roles/arch/systemd/files/user/irc-notifications.service
[Unit]
Description=IRC Client Notifications

[Service]
ExecStart=/home/markl/GitApps/IRCNotify/client.arch.js

[Install]
WantedBy=default.target
